# config/timemachine.service
[Unit]
Description=TimeMachine Backup Service
Documentation=https://github.com/ronaldjonkers/timemachine-backup-linux
After=network.target

[Service]
Type=simple
User=timemachine
Group=timemachine
WorkingDirectory=/opt/timemachine-backup-linux
ExecStart=/opt/timemachine-backup-linux/bin/tmserviced.sh --foreground
ExecReload=/bin/kill -HUP $MAINPID

# Auto-restart: always restart, wait 10s between attempts
Restart=always
RestartSec=10

# Allow up to 5 restarts within 120 seconds before giving up
StartLimitIntervalSec=120
StartLimitBurst=5

# Watchdog: systemd kills the service if it doesn't respond within 120s
WatchdogSec=120

StandardOutput=journal
StandardError=journal
SyslogIdentifier=timemachine

# Runtime directories (systemd creates these automatically with correct ownership)
RuntimeDirectory=timemachine
RuntimeDirectoryMode=0750
StateDirectory=timemachine
LogsDirectory=timemachine

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false

# Environment
EnvironmentFile=-/opt/timemachine-backup-linux/.env

[Install]
WantedBy=multi-user.target
